# launchd/local.codex-sentinel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.codex-sentinel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Codex Sentinel.app/Contents/MacOS/codex-sentinel</string>
    <string>lifecycle</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/gosu/.codex-sentinel/lifecycle.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/gosu/.codex-sentinel/lifecycle.err.log</string>
</dict>
</plist>
